# config/files/cron.d
## COMMANDS
# sudo crontab -l
# sudo crontab -e

# Run each day at 8:05
05 8 * * * sudo systemctl reboot

# Run each day at 17:05
05 17 * * * sudo systemctl reboot

# Run each day at 1:05
05 1 * * * sudo systemctl reboot